# profiles/supervisord/test-profiles.conf
[unix_http_server]
file=/tmp/test-profiles-supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/test-profiles-supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/test-profiles-supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/test-profiles-supervisor.sock ; use a unix:// URL  for a unix socket

[program:webapp]
command=/home/communityprofiles/.virtualenvs/2017Profiles/bin/run-webapp prod-2017.yaml
environment=HOME="/home/communityprofiles",USER="communityprofiles",VIRTUAL_ENV="/home/communityprofiles/.virtualenvs/2017Profiles",LANG=en_US.UTF-8

[eventlistener:crashmail]
command=/home/communityprofiles/.virtualenvs/2017Profiles/bin/crashmail -a -m info@216software.com
events=PROCESS_STATE

[eventlistener:memmon]
command=/home/communityprofiles/.virtualenvs/2017Profiles/bin/memmon -a 1500MB -m info@216software.com
events=TICK_60
